# launchd/com.ai.market-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ai.market-check</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/naokijodan/ai-scripts/run-claude.sh</string>
		<string>--template</string>
		<string>market-check</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/naokijodan/ai-scripts/logs/launchd-market-check.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/naokijodan/ai-scripts/logs/launchd-market-check-error.log</string>
	<key>RunAtLoad</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/naokijodan/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
		<key>HOME</key>
		<string>/Users/naokijodan</string>
	</dict>
</dict>
</plist>
