# launchd/com.n8n.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.n8n.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/naokijodan/.npm-global/bin/n8n</string>
		<string>start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/naokijodan/ai-scripts/logs/n8n.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/naokijodan/ai-scripts/logs/n8n-error.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/naokijodan/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
		<key>HOME</key>
		<string>/Users/naokijodan</string>
		<key>N8N_PORT</key>
		<string>5678</string>
		<key>N8N_PROTOCOL</key>
		<string>http</string>
		<key>N8N_HOST</key>
		<string>localhost</string>
		<key>NODE_FUNCTION_ALLOW_BUILTIN</key>
		<string>child_process,crypto,fs,http,https,os,path,url</string>
		<key>LINE_CHANNEL_ACCESS_TOKEN</key>
		<string>BYx5qBj5aLhHRBfPHkQWLTqDDiL/Ds1s+lKS/f8yWDKDOwfOgC5epMoF346KQ2CkuDjH8vKbCl986VWvsy01jdqKI3p530cXPCXDTYDRH4PHX11vsodbLVZr3+d0hRg55G4mL7vAotNPJhgYZw/CnAdB04t89/1O/w1cDnyilFU=</string>
		<key>LINE_CHANNEL_SECRET</key>
		<string>f46f0d2daf1e79db8f1830b1a9e2884e</string>
	</dict>
</dict>
</plist>
